# root/etc/services.d/homeassistant/run
#!/usr/bin/with-contenv bash

python -m homeassistant --config /config
